# roles/setup_unattended_upgrades/templates/apt-daily.timer.j2
[Timer]
OnCalendar=
RandomizedDelaySec=
OnCalendar=*-*-* 10:00
RandomizedDelaySec=1h