# roles/deploy_mail_server/templates/nextcloud/nextcloud-cron.service.j2
[Unit]
Description=nextcloud-cron
Requires=network.target

[Service]
User=www-data
Type=oneshot
ExecStart=/usr/bin/php -f /var/www/nextcloud/cron.php
